# NMTK_apps/NMTK_server/crons/verify_celery.cron
# Cron job to verify celery is working
# this file should be updated with the appropriate (NMTK_INSTALL_PATH) and copied to
# /etc/cron.d

# This cron job will verify that celery is working properly
BASE_DIR=NMTK_INSTALL_PATH

SHELL=/bin/bash

0 * * * * www-data $BASE_DIR/venv/bin/python $BASE_DIR/NMTK_apps/manage.py verify_celery --email-admins
